# whatsapp-bot.service
[Unit]
Description=WhatsApp AI Business Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0


[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/whatsapp-automation
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/%i/whatsapp-automation/.env

# Run the API server (which includes the WhatsApp bot logic)
ExecStart=/home/%i/whatsapp-automation/venv/bin/python -m uvicorn src.api.webhook:app --host 0.0.0.0 --port 8000

# Auto-restart settings
Restart=always
RestartSec=5
StartLimitBurst=10
StartLimitWindow=60
WatchdogSec=90

# Logging
StandardOutput=append:/home/%i/whatsapp-automation/logs/bot.log
StandardError=append:/home/%i/whatsapp-automation/logs/bot-error.log


[Install]
WantedBy=multi-user.target
